# install/projects_claude-usage-notch-server_ping.timer
[Unit]
Description=Daily Claude ping timer
Requires=projects_claude-usage-notch-server_ping.service

[Timer]
# Change OnCalendar to adjust the schedule (man systemd.time for syntax)
OnCalendar=*-*-* 06:00:00

[Install]
WantedBy=timers.target
