# install/projects_claude-usage-notch-server_scheduler.service
[Unit]
Description=Claude usage notch scheduler (quota poll + daily ping)
After=multi-user.target

[Service]
WorkingDirectory=/home/mnalavadi/claude-usage-notch-server
Type=idle
ExecStart=/home/mnalavadi/.local/bin/uv run scheduler
User=mnalavadi

[Install]
WantedBy=multi-user.target
